# systemd/crankk.timer
[Unit]
Description=Runs Crankk Check every 1 Minute

[Timer]
OnBootSec=15
OnUnitActiveSec=15
AccuracySec=1ms

[Install]
WantedBy=timers.target
