# systemd/crankk.service
[Unit]
Description=Crankk Service

[Service]
Type=oneshot
ExecStart=/etc/monitor-scripts/crankk.sh

[Install]
WantedBy=multi-user.target
