# systemd/myst.service
[Unit]
Description=Mysterium Service

[Service]
Type=oneshot
ExecStart=/etc/monitor-scripts/myst.sh

[Install]
WantedBy=multi-user.target
